# docker/rootfs/etc/services.d/sshd/finish
#!/usr/bin/execlineb -P
s6-svscanctl -t /var/run/s6/services
