# docker/rootfs/etc/services.d/sshd/run
#!/usr/bin/execlineb -P
/bin/ssh-honeypot-wrapper